# systemd/ping_target.service
[Unit]
Description=Ping Target IP

[Service]
Type=oneshot
ExecStart=/usr/local/bin/ping_target.sh
